# awesome_software/chirpstack_concentratord/uci-defaults/91-save-gateway-id
#!/bin/sh
(
    /opt/awesome_linxdot/awesome_software/chirpstack_concentratord/get_gateway_id.sh &
) &
exit 0
